# deploy/denox.service
# systemd unit for DenoX (compiled binary deployment).
# 1. Build:   deno task compile
# 2. Copy:    ./dist/app -> /opt/denox/app  and  .env -> /opt/denox/.env
# 3. Install: sudo cp deploy/denox.service /etc/systemd/system/
#             sudo systemctl daemon-reload && sudo systemctl enable --now denox

[Unit]
Description=DenoX application server
After=network.target

[Service]
Type=simple
User=denox
Group=denox
WorkingDirectory=/opt/denox
EnvironmentFile=/opt/denox/.env
ExecStart=/opt/denox/app
Restart=on-failure
RestartSec=3

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/denox

[Install]
WantedBy=multi-user.target
